# cookbooks/config/files/default/rc.local
#!/bin/sh
/bin/sleep 10
/usr/local/bin/uwsgi --ini /home/ubuntu/django/uwsgi.ini --daemonize /var/log/django.log --touch-reload=/home/ubuntu/reload.ini