# cookbooks/config/files/default/rc.local
#!/bin/sh
sleep 10
/usr/local/bin/uwsgi --ini /home/ubuntu/exercise-8/uwsgi.ini --daemonize /var/log/django-contacts.log